# dist/systemd/infra-api.service
[Unit]
Description=Raku infra-api webhook receiver
After=network.target

[Service]
User=infra-api
Group=infra-api
ExecStart=/var/lib/infra-api/rakudo/bin/raku -I /var/lib/infra-api/__current \
    /var/lib/infra-api/__current/bin/webservice
WorkingDirectory=/var/lib/infra-api/__current
Restart=always

[Install]
WantedBy=multi-user.target
